# holidays/etc/cron.d/holidays
#
# Remove old holiday timeperiods in February
#
1 1	12 2 *	$OMD_ROOT/local/bin/holidays.py delete_old
#
# Automatically add new holiday timeperiods in November
#
1 1	12 11 *	$OMD_ROOT/local/bin/holidays.py add_auto_holidays
